# touchkio-watchdog.timer
[Unit]
Description=Run TouchKIO health watchdog every 2 minutes

[Timer]
OnBootSec=3min
OnUnitActiveSec=2min
AccuracySec=15s
Persistent=false

[Install]
WantedBy=timers.target
